# .config/systemd/user/gallery-dl.timer
[Unit]
Description=Daily Pixiv bookmark downloader

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
