# dot_config/systemd/user/yt-dlp.service.tmpl
[Service]
Type=oneshot

WorkingDirectory=%h/Videos/BiliBili
ExecStart=/usr/bin/yt-dlp {{ .data.yt_dlp.args | join " "}}