# root/dot_config/systemd/user/pixivfe-wrapper.socket
[Socket]
ListenStream=51833

[Install]
WantedBy=sockets.target