# root/dot_config/systemd/user/proxy-to-pixivfe.service
[Unit]
Requires=pixivfe.service
Requires=proxy-to-pixivfe.socket

[Service]
ExecStart=/usr/lib/systemd/systemd-socket-proxyd localhost:43580