# root/dot_config/systemd/user/proxy-to-pixivfe.socket
[Unit]
Description=PixivFE Socket
StopWhenUnneeded=yes

[Socket]
ListenStream=51834

[Install]
WantedBy=sockets.target
